# docker/supervisorConfigs/caeapp.conf
[program:caeapp]
directory=/usr/src/app/cae-app/
command=npm start
autostart=true
autorestart=true
stderr_logfile=/var/log/cae.err.log
stdout_logfile=/var/log/cae.out.log
