# docker/supervisorConfigs/caefrontend.conf
[program:caefrontend]
directory=/usr/src/app/widgets/
command=grunt connect
autostart=true
autorestart=true
stderr_logfile=/var/log/cae.err.log
stdout_logfile=/var/log/cae.out.log
