# docker/supervisorConfigs/syncmeta.ini
[program:syncmeta]
directory=/usr/src/app/syncmeta/widgets
command=grunt connect
autostart=true
autorestart=true
stderr_logfile=/var/log/syncmeta.err.log
stdout_logfile=/var/log/syncmeta.out.log